# init/gosearch.service
[Unit]
Description=gosearch file indexing server

[Service]
Type=simple
ExecStart=/usr/bin/gosearchServer

[Install]
WantedBy=multi-user.target
